# .config/systemd/user/wallpaper-random.timer
[Unit]
Description=Rotate wallpaper every 10 minutes

[Timer]
OnActiveSec=0
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
