# .config/systemd/user/niri.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/sleep 1
